# tools/supervisorapp.conf
[program:blog]
directory=/home/CVAR/study/interest/blogs/github_page/oukohou.github.io/
command=python3 -m server -p 8000 -r /home/server/blog/
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/home/data/logs/supervisord-blog.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stderr_logfile=/home/data/logs/supervisord-blog.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB
user = server
environment = HOME="/home/server",
              USER="server",
              PYTHONPATH="/home/CVAR/study/interest/blogs/github_page/oukohou.github.io",
              PYTHONUNBUFFERED=1
